# packages/es7s/es7s-3.83.0.dev0.tar.gz/es7s-3.83.0.dev0/es7s/data/systemd/40.xbindkeys.service
[Unit]
Description=xbindkeys
#ConditionUser=%UID
#After=multi-user.target
#StartLimitIntervalSec=0

[Service]
Type=simple
KillMode=process
Restart=always
RestartSec=1
ExecStart=/usr/bin/xbindkeys -n -f /home/%USER/.es7s/data/conf/xbindkeysrc --verbose
User=%UID
Environment="DISPLAY=:0"
Environment="ES7S_BIN_PATH=/home/%USER/.es7s/bin"

[Install]
WantedBy=graphical-session.target
